# ddns/ddns.cron
*/5 * * * * root /opt/pi_scripts/ddns/u_v6 &>/dev/null
