# ddns/ddns.cron
*/5 * * * * root /opt/ddns/u_v6 >/dev/null
