# ddns/e_rebooter.cron
23 7 * * * root /opt/ddns/e_rebooter &>/dev/null
