# ddns/i_test.cron
18 1 * * * root /opt/ddns/i_test &>/dev/null
04 7 * * * root /opt/ddns/i_test &>/dev/null
26 13 * * * root /opt/ddns/i_test &>/dev/null
06 19 * * * root /opt/ddns/i_test &>/dev/null
