# ddns/refresh.cron
02 12 * * * root rm -f /tmp/ip6tmp_duck
33 21 * * * root rm -f /tmp/ip6tmp_duck
